# master/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:broker]
command=node /usr/app/broker/index.js
